# atualizacao-cache/crontab
0 1 * * * cd /path-to/atualizacao-cache/ ;  Rscript main.R > /path-to/atualizacao-cache-$(date).log
