# atualizacao-cache/crontab
0 1 * * * Rscript /path/main.R > /path/atualizacao-cache-$(date).log
