# atualizacao-cache/crontab
0 1 * * * /path/Rscript main.R > /path/atualizacao-cache-$(date).log
